# supervisor.ini
[supervisord]
nodaemon = false

[program:grooo]
directory = /usr/src/app
command = gunicorn --bind 0.0.0.0:5000 -k gevent --workers 4 manage:app